# crontab
*/5 * * * * python /invest_dash/backend/crontest.py >> /cronlog.txt